# systemd/user/kmonad.service
[Unit]
Description=KMonad service

[Service]
Type=simple
ExecStart=/bin/kmonad %h/.config/kmonad/kmonad.kbd

[Install]
WantedBy=default.target